# systemd-unit/log-weather.timer
[Unit]
Description=Log current temperature and humidity to DB

[Timer]
OnCalendar=*:0/60

[Install]
WantedBy=timers.target
